# printer3d/octoprint/supervisord_socat.conf

[program:socat]
command=socat -d -d pty,raw,link=/dev/virtualcom0,user=root, tcp:192.168.2.14:23
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
autostart=false
autorestart=false
